# deployment/configs/linux/rc.local
#!/bin/bash
export PATH=$PATH:/root/.local/bin:/root/.local/bin
cd /root/dv0vd.xyz
make on-startup
exit 0
