# deployment/configs/linux/rc.local
#!/bin/bash
export PATH=$PATH:/root/.local/bin
cd /root/dv0vd-backup
make on-startup > /var/log/on-startup.log
exit 0
