# deployment/configs/linux/rc.local
#!/bin/bash
export PATH=$PATH:/root/.local/bin:/root/.local/bin
cd /root/dv0vd.xyz-backup
make logs-clear
make schedule-midnight-reboot
exit 0
